# scripts/send_queued_emails.conf
description "Send emails queued manually in the halsoraport admin"
author "Pedro"
start on runlevel [2345]
respawn

exec /var/www/hrpt/bin/python /var/www/hrpt/manage.py send_queued_emails
